# cookbooks/base/files/etc/cron.d/aws-mon-linux
*/5 * * * * root /usr/local/bin/aws-mon.sh --all-items --disk-path=/ --from-cron
